# cron/crontab
30 4 * * * /bin/bash -l -c 'cd /usr/src/app && RAILS_ENV=production bundle exec rake import:data_dump' > /proc/1/fd/1

* * * * * /bin/bash -l -c 'echo test5' > /proc/1/fd/1
